# openwrt/files/etc/init.d/mihomo
#!/bin/sh /etc/rc.common
START=94
STOP=10
USE_PROCD=1
log_level_num() {
	case "$(uci -q get purewrt.settings.log_level)" in
		error|err) echo 0 ;;
		info|notice) echo 2 ;;
		debug) echo 3 ;;
		*) echo 1 ;;
	esac
}
log_msg() {
	level="$1"
	shift
	case "$level" in
		error|err) need=0 ;;
		info|notice) need=2 ;;
		debug) need=3 ;;
		*) need=1 ;;
	esac
	[ "$(log_level_num)" -ge "$need" ] && logger -t mihomo "$level $*"
}
start_service() {
	local bin workdir config
	config_load purewrt
	config_get bin settings mihomo_bin '/usr/bin/mihomo'
	config_get workdir settings workdir '/etc/purewrt'
	config_get config settings mihomo_config '/etc/purewrt/generated/mihomo.yaml'
	log_msg warn "start: mihomo service start requested"
	log_msg info "start: bin=$bin workdir=$workdir config=$config"
	mkdir -p "$(dirname "$config")"
	if [ ! -s "$config" ]; then
		log_msg info "start: mihomo config missing, generating PureWRT config"
		/usr/bin/purewrt generate 2>&1 | logger -t mihomo-purewrt-generate
	fi
	# Enable TCP Fast Open (client+server) so mihomo's tfo: outbounds can
	# negotiate it. Best-effort; harmless if the knob is absent.
	sysctl -w net.ipv4.tcp_fastopen=3 >/dev/null 2>&1
	procd_open_instance
    # Wrap the exec to (1) raise RLIMIT_NOFILE and (2) pin GOMAXPROCS to all
    # cores. nofile: a proxy core holds a socket + pipes per live connection and
    # health-checks every node in its url-test/load-balance groups; the default
    # (~4096) is exhausted fast on a real subscription → "socket: too many open
    # files". `procd_set_param limits nofile` is clamped by procd's (pid 1) own
    # hard limit (4096), so it does NOT work — a root child instead raises its
    # own limit up to fs.nr_open via ulimit. GOMAXPROCS: Go already defaults to
    # NumCPU, but set it explicitly so a future cgroup CPU quota can't silently
    # throttle the proxy to one core.
    procd_set_param command /bin/sh -c 'ulimit -n 1000000 2>/dev/null; c=$(grep -c ^processor /proc/cpuinfo 2>/dev/null); [ "${c:-0}" -gt 1 ] && export GOMAXPROCS="$c"; exec "$1" -d "$2" -f "$3"' mihomo "$bin" "$workdir" "$config"
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn
    procd_close_instance
	log_msg info "start: mihomo procd instance configured"
}
reload_service() { log_msg info "reload: sending signal to mihomo"; procd_send_signal mihomo; }
